# system-scripts/nextcloud-db-dump/nextcloud-db-dump.service
[Unit]
Description=Consistent logical dump of the Nextcloud database (for offsite backup)
Documentation=file:/opt/system-scripts/nextcloud-db-dump/README.adoc
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/opt/system-scripts/nextcloud-db-dump/nextcloud-db-dump
# Harmless if it overruns; the atomic rename means a slow dump never yields a
# torn file. Cap it so a hung DB does not wedge the unit forever.
TimeoutStartSec=1800
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=6

[Install]
WantedBy=multi-user.target
