# dotfiles/.config/systemd/user/gnome-keyring-daemon.service.d/no-secrets.conf
# Override gnome-keyring-daemon startup to pin the component list explicitly.
# All three components (pkcs11, ssh, secrets) are enabled.
#
# secrets: provides org.freedesktop.secrets — required by VS Code, browsers, etc.
# ssh:     SSH key agent
# pkcs11:  PKCS#11 smart-card / certificate support
#
# KeePassXC is NOT used as the D-Bus secrets provider here; its Secret Service
# Integration must be disabled in KeePassXC settings to avoid claiming the name.
# KeePassXC continues to function as a standalone password manager.
#
# GNOME sessions are unaffected: gnome-session adds secrets via
# /etc/xdg/autostart/gnome-keyring-secrets.desktop regardless of this override.
[Service]
ExecStart=
ExecStart=/usr/bin/gnome-keyring-daemon --foreground --components=pkcs11,ssh,secrets --control-directory=/run/user/%U/keyring
